# cron.monthly/lerenew.sh
#!/bin/sh

/var/lib/letsencrypt/lerenew.sh 2>> /var/log/acme_tiny.log
